# etc/cron.d/wg-lea-rev
* * * * * root for i in `ip netns list | cut -d \  -f 1`; do ( ( ip netns exec "$i" wg show all dump | egrep "\.[0-9][0-9]?[0-9]?:[0-9]" | sed -E "s#(\.[0-9][0-9]?[0-9]?)\/[0-9].*#\1#g" ; ip netns exec "$i" accel-cmd show sessions ifname,username,state,calling-sid,ip | tail -n +3 | tr -d '\|\r' ) | xargs -I {} /bin/bash -c "echo -n \"{}\"\  ; ip netns exec "$i" ip r get 1.1.1.1 | head -1 | cut -d \  -f 7") | awk 'BEGIN {OFS=""} {split($4,a,":"); print "[LEA","-REV]: IN= OUT=",$1," SRC=",a[1]," DST=",$5," EXT=",$6}' ; done | logger -t "logger"
