# etc/cron.d/wg-lea-rev
* * * * * root /wg-lea-rev.sh | logger -t "logger"
