# etc/systemd/system/wg-quick-ns.target
[Unit]
Description=WireGuard Tunnels via wg-quick-ns(8)
